# systemd/gptbot-backup.timer
[Unit]
Description=Run GPTBOT daily backup

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=gptbot-backup.service

[Install]
WantedBy=timers.target
